# kata-guest-base/extra/etc/systemd/system/kata-agent.service.d/10-ghcr-auth.conf
# Point CDH (spawned by kata-agent) at the in-guest GHCR auth file that
# systemd-tmpfiles copies from the dm-verity-baked
# /etc/c8s/ghcr-auth.json. Without this, image-rs falls back to
# Anonymous and every private ghcr pull returns 401.
#
# Kata 3.30's kata-agent parses `agent.image_registry_auth=…` off the
# kernel cmdline for a subset of paths — but the guest-pull code path
# through image-rs does NOT honour that arg unless CDH reads the URI
# from its config. CDH reads
# `CDH_DEFAULT_IMAGE_AUTHENTICATED_REGISTRY_CREDENTIALS` from the
# environment as a documented fallback (see
# confidential-data-hub/hub/src/config.rs).
[Service]
Environment=CDH_DEFAULT_IMAGE_AUTHENTICATED_REGISTRY_CREDENTIALS=file:///run/image-security/auth.json
